# tronWallet/transactions_demon/config/supervisord.conf
[supervisord]
user=
nodaemon=

[program:eth-daemon]
command=
user=root
numprocs=
startsecs=
autostart=
autorestart=
process_name=