# solanaWallet/api/config/supervisord.conf
[supervisord]
user=
nodaemon=

[program:eth-api]
command=
user=
numprocs=
startsecs=
autostart=
autorestart=
process_name=